# private_dot_config/systemd/user/mako.service.d/override.conf
[Service]
ExecCondition=/bin/sh -c '\
[ -n "$SWAYSOCK" ] || \
[ -n "$NIRI_SOCKET" ] || \
[ "$XDG_CURRENT_DESKTOP" = "sway" ] || \
[ "$XDG_CURRENT_DESKTOP" = "niri" ]'
